# netPi/etc/systemd/system/airupnp.timer
[Unit]
Description=Runs Airconnect after 2 minutes

[Timer]
OnBootSec=2min
Unit=airupnp.service

[Install]
WantedBy=multi-user.target
